# service/xautolock/run
#!/usr/bin/bash
exec xautolock -time 5 -locker "i3lock -n -c fbf1c7 -i ${HOME}/.config/lock.png"
