# .local/service/lanta/run
#!/usr/bin/bash
cd ${HOME}
exec ${HOME}/.cargo/bin/lanta
