# service/cnx/run
#!/usr/bin/bash
exec ${HOME}/.cargo/bin/cnx
