# .local/service/unclutter/run
#!/usr/bin/bash

exec unclutter -idle 1
