# service/polybar/run
#!/usr/bin/bash
export MONITOR=DP-1
exec polybar main 2>&1
